# deploy/systemd/vela-bundler-alert.service
[Unit]
Description=Vela Bundler failure alert — pushes a Telegram message when the bundler unit FAILS
# Fired by vela-bundler.service's OnFailure=. This is the layer that reports "the bundler is
# DOWN AND STAYED DOWN" (start-limit exhaustion included) — the process's own alerting cannot
# report its own death. Requires TELEGRAM_BOT_TOKEN + TELEGRAM_CHAT_ID in vela.env.

[Service]
Type=oneshot
EnvironmentFile=/opt/vela-bundler/data/vela.env
# %H = hostname. --max-time keeps a Telegram outage from wedging the oneshot.
ExecStart=/usr/bin/curl -s --max-time 10 \
  --data-urlencode "chat_id=${TELEGRAM_CHAT_ID}" \
  --data-urlencode "text=🔴 Vela Bundler systemd unit FAILED on %H — the bundler is DOWN (restart limit likely exhausted). Every user op is failing until it is restarted. journalctl -u vela-bundler for the cause." \
  "https://api.telegram.org/bot${TELEGRAM_BOT_TOKEN}/sendMessage"
